# config/systemd/user/yandex-disk.service
[Unit]
Description=yandex-disk 
After=network.target

[Service]
ExecStart=yandex-disk start
